# docker/pgsql.ini
[program:pgsql]
startsecs=10
autostart=true
autorestart=true
priority=100
command=/usr/local/bin/docker-entrypoint.sh postgres -c 'config_file=/etc/postgresql/postgresql.conf'
