# docker/jobsrv/supervisor/jc.ini
[program:jc]
startsecs=10
autostart=true
autorestart=true
priority=300
command=/root/jc_init.sh
