# docker/jobsrv/supervisor/rpcswitch.ini
[program:rpcswitch]
startsecs=10
autostart=true
autorestart=true
priority=300
command=/root/rpcswitch_init.sh
